# service/relayctl-scheduler.service
[Unit]
Description=relayctl scheduler service

[Service]
Type=simple
ExecStart=/usr/bin/python /etc/relayctl/service/main.py